# dashboard/src/lib/mock/investigations.json
{
  "cases": [
    {
      "id": "INV-2026-001",
      "title": "Lateral Movement — U4521@DOM2 to Domain Controller",
      "status": "In Progress",
      "severity": 4,
      "created": "2026-02-10T14:32:00Z",
      "updated": "2026-02-11T09:15:00Z",
      "assignee": "Nethra",
      "eventCount": 47,
      "description": "Anomalous Kerberos authentication chain detected: U4521@DOM2 authenticated from C102 to C4501, then C4501 to C892, then C892 to DC01 (domain controller). Pattern matches MITRE T1021 lateral movement via Remote Services.",
      "tags": ["lateral-movement", "T1021", "kerberos", "domain-controller"],
      "hosts": ["C102", "C4501", "C892", "DC01"],
      "users": ["U4521@DOM2"]
    },
    {
      "id": "INV-2026-002",
      "title": "Suspicious PowerShell Execution Chain",
      "status": "Open",
      "severity": 3,
      "created": "2026-02-11T02:18:00Z",
      "updated": "2026-02-11T02:18:00Z",
      "assignee": "Unassigned",
      "eventCount": 12,
      "description": "PowerShell process spawned by cmd.exe with encoded command argument on host C3847. Process tree shows unusual parent chain: svchost.exe → cmd.exe → powershell.exe with -EncodedCommand flag.",
      "tags": ["execution", "T1059", "powershell", "encoded-command"],
      "hosts": ["C3847"],
      "users": ["U8921@DOM1"]
    },
    {
      "id": "INV-2026-003",
      "title": "Data Exfiltration via DNS Tunneling",
      "status": "Open",
      "severity": 4,
      "created": "2026-02-10T22:45:00Z",
      "updated": "2026-02-11T08:30:00Z",
      "assignee": "Marcus Webb",
      "eventCount": 234,
      "description": "Abnormal DNS query volume from C1923: 2,400+ queries in 10 minutes to subdomain pattern matching base64-encoded data tunneling. Destination resolves to external IP outside corporate range.",
      "tags": ["exfiltration", "T1041", "dns-tunneling", "command-and-control"],
      "hosts": ["C1923"],
      "users": ["U3102@DOM3"]
    },
    {
      "id": "INV-2026-004",
      "title": "Credential Dumping — Mimikatz Detection",
      "status": "Closed",
      "severity": 4,
      "created": "2026-02-09T16:20:00Z",
      "updated": "2026-02-10T11:00:00Z",
      "assignee": "Nethra",
      "eventCount": 8,
      "description": "Process execution of mimikatz.exe detected on C587. Binary hash matches known Mimikatz variant. Process accessed LSASS memory (pid 672). Credentials may have been extracted.",
      "tags": ["credential-access", "T1003", "mimikatz", "lsass"],
      "hosts": ["C587"],
      "users": ["U1205@DOM1"]
    },
    {
      "id": "INV-2026-005",
      "title": "Brute Force Authentication — Multiple Failures",
      "status": "Closed",
      "severity": 2,
      "created": "2026-02-09T08:10:00Z",
      "updated": "2026-02-09T14:30:00Z",
      "assignee": "James Morton",
      "eventCount": 1842,
      "description": "1,842 failed NTLM authentication attempts from C9234 targeting 15 different user accounts over 2-hour window. Pattern consistent with password spraying. No successful authentications detected.",
      "tags": ["credential-access", "T1110", "brute-force", "password-spray"],
      "hosts": ["C9234"],
      "users": ["Multiple"]
    },
    {
      "id": "INV-2026-006",
      "title": "Persistence — Scheduled Task Created",
      "status": "In Progress",
      "severity": 3,
      "created": "2026-02-11T06:45:00Z",
      "updated": "2026-02-11T10:20:00Z",
      "assignee": "James Morton",
      "eventCount": 5,
      "description": "New scheduled task 'WindowsUpdateCheck' created on C4102 by U7823@DOM2 executing payload from C:\\Windows\\Temp\\svc.exe. Binary not in approved software inventory. Task set to run every 15 minutes.",
      "tags": ["persistence", "T1053", "scheduled-task"],
      "hosts": ["C4102"],
      "users": ["U7823@DOM2"]
    },
    {
      "id": "INV-2026-007",
      "title": "Network Scan — Internal Reconnaissance",
      "status": "Open",
      "severity": 2,
      "created": "2026-02-11T04:12:00Z",
      "updated": "2026-02-11T04:12:00Z",
      "assignee": "Unassigned",
      "eventCount": 3200,
      "description": "Host C2891 initiated connections to 3,200 unique internal IPs on port 445 (SMB) within 5 minutes. Pattern consistent with network enumeration for file shares. Source process: nmap.",
      "tags": ["discovery", "T1018", "network-scan", "smb"],
      "hosts": ["C2891"],
      "users": ["U5431@DOM1"]
    },
    {
      "id": "INV-2026-008",
      "title": "Privilege Escalation — Service Account Misuse",
      "status": "In Progress",
      "severity": 3,
      "created": "2026-02-10T19:30:00Z",
      "updated": "2026-02-11T07:45:00Z",
      "assignee": "Nethra",
      "eventCount": 23,
      "description": "Service account SYSTEM on C6721 used for interactive logon — normally only batch/service logon types. Subsequent process execution includes cmd.exe and whoami.exe, suggesting hands-on-keyboard activity.",
      "tags": ["privilege-escalation", "T1078", "service-account"],
      "hosts": ["C6721"],
      "users": ["SYSTEM"]
    }
  ]
}
